# raspberry-pi/vitrine-display.service
[Unit]
Description=Vitrine.io E-Ink Display Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/vitrine-display
EnvironmentFile=/home/pi/vitrine-display/.env
ExecStart=/usr/bin/python3 /home/pi/vitrine-display/vitrine_display.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
